# supervisor.ini
[program:percinmedia-telegram]
process_name = %(process_num)02d
command = /usr/bin/php /home/percinmedia/htdocs/percinmedia.net/artisan queue:work --queue=telegram --tries=3 --timeout=3600
autostart = true
autorestart = true
stopasgroup = true
killasgroup = true
user = percinmedia
numprocs = 1
redirect_stderr = true
stdout_logfile = /home/percinmedia/logs/supervisor/queue-telegram.log
stdout_logfile_maxbytes = 5MB
stdout_logfile_backups = 5
